# rag-api.service
[Unit]
Description=RAG Report Generator API
After=network.target

[Service]
Type=notify
User=work
Group=work
WorkingDirectory=/home/work/rag/Project/rag-report-generator
Environment="PATH=/home/work/rag/Project/rag-report-generator/.venv/bin:/usr/local/bin:/usr/bin:/bin"
EnvironmentFile=/home/work/rag/Project/rag-report-generator/.env

ExecStart=/home/work/rag/Project/rag-report-generator/.venv/bin/gunicorn \
    --config /home/work/rag/Project/rag-report-generator/gunicorn.conf.py \
    app.api.main:app

# 재시작 정책
Restart=always
RestartSec=10
StartLimitBurst=5
StartLimitInterval=60

# 보안 설정
NoNewPrivileges=true
PrivateTmp=true

# 로깅
StandardOutput=journal
StandardError=journal
SyslogIdentifier=rag-api

[Install]
WantedBy=multi-user.target
